# mailjaeger.service
[Unit]
Description=MailJaeger - Local AI Email Processing System
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/MailJaeger
Environment="PATH=/home/pi/MailJaeger/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/MailJaeger/venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=4G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/MailJaeger

[Install]
WantedBy=multi-user.target
